# cron.cron
*/1 * * * * export DISPLAY=:0.0 && /home/marcin/PycharmProjects/so/pencrypt.py -a
